# supervisord.conf
[supervisord]
nodaemon=true

[program:c2server]
command=python /app/server/main.py
autostart=true
autorestart=true

[program:implant]
command=/usr/local/bin/implant
autostart=true
autorestart=true
